# systemd/embedded-config-manager.service
[Unit]
Description=Embedded Config Manager agent
Documentation=https://github.com/aykutsp/embedded-config-manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ecm
Group=ecm
Environment=ECM_ROOT=/opt/embedded-config-manager
Environment=ECM_DATA_DIR=/var/lib/embedded-config-manager
Environment=ECM_DRY_RUN=0
ExecStart=/opt/embedded-config-manager/.venv/bin/ecm-agent
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/embedded-config-manager /etc/telemetry /etc/network/interfaces.d
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
